# ansible/playbooks/roles/node/files/e3/harmony.service
[Unit]
Description=harmony service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
WorkingDirectory=/root
ExecStart=/root/node.sh -N staking -1 -P -p /root/.hmy/blskeys/bls.pass -M -D -T explorer -i 3 -A

[Install]
WantedBy=multi-user.target
