# tf/do/files/harmony.service
[Unit]
Description=Harmony Service
Wants=network.target network-online.target
After=network.target network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root
ExecStart=bash /root/node.sh -1 -S -p /root/bls.pass -k /root/bls.key -P
Restart=always

[Install]
WantedBy=multi-user.target 