# ansible/playbooks/roles/node/files/s2/crontab
SHELL=/usr/bin/bash
HOME=/root

40 5,17 * * * systemctl stop harmony; sleep 3; rclone sync $HOME/harmony_db_2 ostn:pub.harmony.one/ostn/harmony_db_2; systemctl start harmony
