# tf/gcp/files/crontab
SHELL=/usr/bin/bash
HOME=/home/gce-user

30 0 * * *  rm -f latest/*.gz
