# tf/gcp/files/harmony-mk.service
[Unit]
Description=Harmony Service
Wants=network.target network-online.target
After=network.target network-online.target

[Service]
Type=simple
User=gce-user
WorkingDirectory=/home/gce-user
ExecStart=bash /home/gce-user/node.sh -1 -S -p /home/gce-user/bls.pass -M -P -D
Restart=always
StandardError=file:/home/gce-user/harmony.err
LimitNOFILE=65536
LimitNPROC=65536

[Install]
WantedBy=multi-user.target 
