# tf/aws.1/files/crontab
30 0 * * *  bash -c "sleep $[RANDOM % 90]m ; /home/ec2-user/uploadlog.sh"
