# tf/aws.1/files/harmony-mk.service
[Unit]
Description=harmony service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=ec2-user
WorkingDirectory=/home/ec2-user
ExecStart=/home/ec2-user/node.sh -1 -S -P -p /home/ec2-user/bls.pass -M -D
StandardError=file:/home/ec2-user/harmony.err
LimitNOFILE=65536
LimitNPROC=65536

[Install]
WantedBy=multi-user.target
