# ansible/playbooks/roles/node/files/s1/crontab
SHELL=/usr/bin/bash
HOME=/root

35 5,17 * * * systemctl stop harmony; sleep 3; rclone sync $HOME/harmony_db_1 ostn:pub.harmony.one/ostn/harmony_db_1; systemctl start harmony
