# ansible/playbooks/roles/node/files/crontab
SHELL=/usr/bin/bash
HOME=/root

30 5,17 * * * systemctl stop harmony; sleep 3; rclone sync $HOME/harmony_db_0 mainnet:pub.harmony.one/ostn/harmony_db_0; systemctl start harmony
