# ansible/playbooks/roles/explorer_fe/templates/harmony-explorer-fe.service.j2
[Unit]
Description=harmony explorer frontend service
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=1
User={{ USER.user }}
Environment=SSL_CERT_FILE=/etc/letsencrypt/live/{{ domain }}/fullchain.pem
Environment=SSL_KEY_FILE=/etc/letsencrypt/live/{{ domain }}/privkey.pem
WorkingDirectory={{ USER.home }}/harmony-dashboard
ExecStart=/bin/bash -l -c 'cd {{ USER.home }}/harmony-dashboard && npm run start'
StandardError=syslog
SyslogIdentifier=harmony-explorer-fe
StartLimitInterval=0
LimitNOFILE=65536
LimitNPROC=65536

[Install]
WantedBy=multi-user.target