# ansible/playbooks/roles/node/templates/harmony.service.j2
[Unit]
Description=harmony service
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=1
User={{ USER.user }}
WorkingDirectory={{ USER.home }}
{% if node_type_var == "Explorer" %}
ExecStart={{ USER.home }}/node.sh -N {{ network }} -T explorer -i {{ shard }} -A1MSPD
{% else %}
ExecStart={{ USER.home }}/node.sh -N {{ network }} -1MSPD
{% endif %}
StandardError=syslog
SyslogIdentifier=harmony
StartLimitInterval=0
LimitNOFILE=65536
LimitNPROC=65536

[Install]
WantedBy=multi-user.target
