# ansible/playbooks/roles/node/templates/harmony.service.j2
[Unit]
Description=harmony service
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=1
User={{ USER.user }}
WorkingDirectory={{ USER.home }}
{% if 'exp' in inventory %}
ExecStart={{ USER.home }}/node.sh -N {{ network }} -T explorer -i {{ shard }} -A -S -M -D -P -1
{% else %}
ExecStart={{ USER.home }}/node.sh -N {{ network }} -S -M -D -P -1
{% endif %}
StandardError=syslog
SyslogIdentifier=harmony
StartLimitInterval=0
LimitNOFILE=65536
LimitNPROC=65536

[Install]
WantedBy=multi-user.target
