# tools/harmony.service.template
[Unit]
Description=harmony service
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=1
User=%%USER%%
WorkingDirectory=%%HOME%% 
ExecStart=%%HOME%%/node.sh -1 -S -P -M -D -N %%NETWORK%% -p %%HOME%%/%%BLSPASS%% -T %%NODETYPE%% %%EXTRA%%
StandardError=syslog
SyslogIdentifier=harmony
StartLimitInterval=0
LimitNOFILE=65536
LimitNPROC=65536

[Install]
WantedBy=multi-user.target
